# AV-postfix-2.timer
[Unit]
Description=AV-postfix-1.timer

[Timer]
OnBootSec=40sec
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target

